# deploy/templates/nexora-node-agent.service
[Unit]
Description=Nexora Node Agent
After=network.target

[Service]
Type=simple
User=nexora
Group=nexora
Environment=NEXORA_NODE_AGENT_HOST=127.0.0.1
Environment=NEXORA_NODE_AGENT_PORT=38121
Environment=NEXORA_STATE_PATH=/opt/nexora/var/state.json
Environment=NEXORA_API_TOKEN_FILE=/opt/nexora/var/api-token
ExecStart=/opt/nexora/venv/bin/nexora-node-agent
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/nexora/var /tmp/nexora-export
ProtectHome=true

[Install]
WantedBy=multi-user.target
