# ynh-package/conf/systemd.service
[Unit]
Description=Nexora Control Plane (YunoHost)
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__
Environment=NEXORA_CONTROL_PLANE_HOST=127.0.0.1
Environment=NEXORA_CONTROL_PLANE_PORT=__PORT__
Environment=NEXORA_STATE_PATH=__DATA_DIR__/state.json
Environment=NEXORA_API_TOKEN_FILE=__DATA_DIR__/api-token
Environment=NEXORA_API_TOKEN_ROLE_FILE=/etc/nexora/api-token-roles.json
Environment=NEXORA_OPERATOR_ONLY_ENFORCE=1
Environment=NEXORA_DEPLOYMENT_SCOPE=operator
Environment=NEXORA_REPO_ROOT=__INSTALL_DIR__/repo
ExecStart=__INSTALL_DIR__/venv/bin/nexora-control-plane
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__DATA_DIR__ /var/log/__APP__ /tmp/nexora-export
ProtectHome=true

[Install]
WantedBy=multi-user.target
